# infra/systemd/bridge.service
[Unit]
Description=Phoenix Bridge FastAPI Service
Documentation=https://github.com/phoenix-trading/openclaw
After=network.target openclaw.service
Requires=openclaw.service

[Service]
Type=simple
User=phoenix
Group=phoenix
WorkingDirectory=/opt/phoenix/openclaw/bridge
Environment=PYTHONPATH=/opt/phoenix/openclaw/bridge
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/uvicorn src.main:app --host 0.0.0.0 --port 18800
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=phoenix-bridge

[Install]
WantedBy=multi-user.target
