# infra/systemd/openclaw.service
[Unit]
Description=Phoenix OpenClaw Trading Process
Documentation=https://github.com/phoenix-trading/openclaw
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=phoenix
Group=phoenix
WorkingDirectory=/opt/phoenix/openclaw
ExecStart=/usr/local/bin/openclaw
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
